# agent-home/deploy/agent-home.service
# systemd unit for agent-home (FG-20 mobile PWA, Next.js BFF) on the prod box.
#
# Serves the production build on loopback only (127.0.0.1:3100, set via PORT in
# the EnvironmentFile; :3000 is the WhatsApp bridge). Caddy terminates TLS on
# the public subdomain and reverse-proxies to it (see Caddyfile.agent-home). The
# Next server makes server-side calls to the Python AI layer at
# AGENT_HOME_API_URL (127.0.0.1:9119) — the Python API is never exposed on the
# public agent-home origin.
#
# Install (adjust paths/User to the box; see deploy/DEPLOY.md):
#   sudo cp agent-home/deploy/agent-home.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now agent-home.service
#
# Secrets are read from an EnvironmentFile that is NOT committed
# (agent-home/agent-home.env — NOT .env.production, which Next's dotenv would
# auto-load and choke on; 0600).

[Unit]
Description=agent-home (FG-20 mobile PWA, Next.js) behind Caddy
Documentation=https://github.com/leolau/ai-prentice-4-all/blob/develop/agent-home/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/data/agent-home-app/agent-home
EnvironmentFile=/opt/data/agent-home-app/agent-home/agent-home.env
# PATH must resolve node/npm for the service user. Verify with `command -v npm`
# on the box during deploy and edit this line if node lives elsewhere (e.g. nvm).
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/data/agent-home-app/agent-home/deploy/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
